# sysdir/etc/systemd/system/smart-sleep.service
[Unit]
Description=Check for idle and sleep

[Service]
ExecStart=/usr/local/bin/smart-sleep.sh
