# sysdir/etc/systemd/system/smart-sleep.timer
[Unit]
Description=Run smart-sleep every minute

[Timer]
OnCalendar=*:*:00
Unit=smart-sleep.service

[Install]
WantedBy=timers.target
